# scripts/com.anchit.portfolio.autosync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.anchit.portfolio.autosync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/anchittandon/Library/Mobile Documents/com~apple~CloudDocs/ANCHIT'S AI HUSTLE/Anchit-Work-Portfolio/scripts/gh-sync.sh</string>
  </array>

  <!-- Run at login and then every 10 minutes (600s). Change as you like. -->
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>600</integer>

  <key>StandardOutPath</key>
  <string>/tmp/anchit-portfolio-autosync.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/anchit-portfolio-autosync.err</string>
</dict>
</plist>
